# deploy/systemd/souqmatbakh-backup.timer
[Unit]
Description=Daily backup timer for KitchenTech database and uploads
Requires=souqmatbakh-backup.service

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
