# deploy/systemd/souqmatbakh-healthcheck.timer
[Unit]
Description=KitchenTech health check timer (every 5 minutes)
Requires=souqmatbakh-healthcheck.service

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true
AccuracySec=30s

[Install]
WantedBy=timers.target
